# myproject/conf/vagrant/upstart.conf
description "myproject - gunicorn server"
start on vagrant-mounted
stop on runlevel [06]
respawn
respawn limit 10 5
exec /vagrant/myproject/conf/vagrant/gunicorn.sh